# legacy/guardian_arquitectura.timer
[Unit]
Description=Guardian Arquitectura Timer
After=network.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
EOF < /dev/null